# launchd/com.claudespeak.report.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.claudespeak.report</string>

  <!-- install.sh substitutes @HOME@ and @PREFIX@ at install time. -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>claudespeak-analyze --days 7 --save</string>
  </array>

  <!-- Daily 9:00 AM local time. Edit to taste. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>@HOME@/.local/state/claudespeak/report.log</string>

  <key>StandardErrorPath</key>
  <string>@HOME@/.local/state/claudespeak/report.err</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
